# deploy/server/cryptosugar-backup.cron
17 3 * * * root /usr/local/sbin/cryptosugar-backup >/var/log/cryptosugar-backup.log 2>&1
17 4 * * 0 root /usr/local/sbin/cryptosugar-verify-backup >/var/log/cryptosugar-backup-verify.log 2>&1
